# config/etapi.conf
[program:etapi]
command = /home/www/etapi/config/gunicorn_start.sh
user = root
stdout_logfile = /home/www/etapi/logs/gunicorn_supervisor.log
redirect_stderr = true
